# lightningos-light/templates/systemd/lightningos-terminal.service
[Unit]
Description=LightningOS Terminal (GoTTY)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=lightningos
Group=lightningos
EnvironmentFile=-/etc/lightningos/secrets.env
ExecCondition=/bin/bash -lc '[ "${TERMINAL_ENABLED}" = "1" ] && [ -n "${TERMINAL_CREDENTIAL}" ]'
ExecStart=/usr/local/sbin/lightningos-terminal
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
